# systemd/hangermon.service
[Unit]
Description=HangerMon human detection service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__HANGERMON_USER__
Group=__HANGERMON_USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=-/etc/default/hangermon
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/python3 -c "import time; time.sleep(2)"
ExecStart=__INSTALL_DIR__/start.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
